# backup/tar.include.list
#/etc/unbound/unbound.conf
#/usr/local/etc/raddb/sites-enabled/dhcp
/etc/sysconfig/network-scripts
/var/www/mikbill
/etc/rc.d/rc.local
/etc/my.cnf
/etc/nginx/conf.d
/etc/crontab
/etc/cron.d
/etc/resolv.conf
/etc/sysctl.conf
/home/backup/SYS_backup.conf
